# sys-auth/seatd/files/seatd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
pidfile=${pidfile:-"/run/${RC_SVCNAME}.pid"}
command="seatd"
command_args="-g video"
